# deploy/tablero-entrenador.service
[Unit]
Description=Tablero del Entrenador (FastAPI/Uvicorn)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/tablero_entrenador
EnvironmentFile=/home/ubuntu/tablero_entrenador/.env
# --proxy-headers + --forwarded-allow-ips: nginx termina TLS y proxya HTTP; así
# FastAPI ve el esquema https correcto (necesario para la cookie de sesión Secure).
ExecStart=/home/ubuntu/tablero_entrenador/.venv/bin/uvicorn server.main:app \
    --host 127.0.0.1 --port 8091 \
    --proxy-headers --forwarded-allow-ips 127.0.0.1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
